# files/etc/systemd/system/lego-renewal.service.j2
[Unit]
Description="Run lego certificate renewals"

[Service]
Type=oneshot
ExecStart={{ lego_path }}/renewals/renewals.sh