# files/etc/systemd/system/lego-renewal.timer.j2
[Unit]
Description=Run lego-renewal.service at specific time

[Timer]
OnCalendar={{ renewal_time }}
Unit=lego-renewal.service
RandomizedDelay={{ randomized_delay }}